# rootfs/etc/inittab
::sysinit:/etc/rc.init
::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc.shutdown
