# rootfs/etc/rc.local
#!/bin/sh

ifconfig lo 127.0.0.1 netmask 255.0.0.0 up

for iface in $(ls /sys/class/net 2>/dev/null | grep -v lo); do
	ip link set "$iface" up
	udhcpc -i "$iface" -n -q -s /usr/share/udhcpc/default.script \
		-p /run/udhcpc-"$iface".pid 2>/dev/null && break
done
